# deploy/mvp-backend-fixed.service
[Unit]
Description=MVP Cost Backend (FastAPI)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/mvp_cost/backend
Environment="PATH=/opt/mvp_cost/venv/bin"
ExecStart=/opt/mvp_cost/venv/bin/python run_prod.py
Restart=always
RestartSec=10
EnvironmentFile=/opt/mvp_cost/.env
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
